# conf/systemd/videorenamer.service
[Unit]
Description=VideoRenamer

[Service]
Type=oneshot
ExecStart=/var/packages/VideoRenamer/target/scripts/videorenamer.sh
